# src/server-daemon/ocland.conf.in
# ocland - OpenCL Land
#
# This configuration file is used for upstart based
# distributions.

description	"Ocland server"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
normal exit 0

script
	OCLAND_SERVER_BIN=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/ocland_server
	OCLAND_SERVER_LOG=/var/log/ocland.log
	exec $OCLAND_SERVER_BIN --log-file $OCLAND_SERVER_LOG &
end script
